# DO4_LinuxMonitoring_v2.0-0-develop/src/09/my_own_node_exporter.service
[Unit]
Description=Basic system monitoring exporter for Prometheus

[Service]
User=root
WorkingDirectory=/media/common_files/Mon2/09
ExecStart=/bin/bash main.sh
Restart=always

[Install]
WantedBy=multi-user.target